# docker/app/production/etc/service/php-fpm/run
#!/bin/sh

php-fpm --nodaemonize --allow-to-run-as-root
